# config/launchd/com.jwalinshah.secret-cache-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jwalinshah.secret-cache-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jwalinshah/.local/bin/secret-cache</string>
        <string>refresh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jwalinshah/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/jwalinshah</string>
    </dict>
    <key>StartInterval</key>
    <integer>86400</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jwalinshah/.cache/quota-core/secret-cache-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jwalinshah/.cache/quota-core/secret-cache-refresh.log</string>
</dict>
</plist>
